# services/agents/install/agentd-cleanup.timer.tmpl
[Unit]
Description=Run Cadence Agent Worktree Cleanup hourly
Requires=agentd-cleanup.service

[Timer]
# Run once per hour. Persistent=true means a missed run (e.g. system was off)
# fires immediately on next boot. RandomizedDelaySec avoids thundering herd
# on multi-host deployments. Matches launchd StartInterval=3600.
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
